# scripts/com.sharkwatch.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sharkwatch.weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/steveelliott/.nvm/versions/node/v24.14.0/bin/node</string>
        <string>/Users/steveelliott/sharks-cron/weekly-report.js</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/steveelliott/sharks-cron/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/steveelliott/sharks-cron/weekly.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/steveelliott</string>
    </dict>
</dict>
</plist>
